# crontab
0 20 * * 2,5 /usr/app/run.sh >> /var/log/cron.log
# Update TWSE/TPEX listings at 08:30 on weekdays
30 8 * * 1-5 python /usr/app/update_listings.py >> /var/log/cron.log
