# crontab
0 20 * * 2,5 /usr/app/run.sh >> /var/log/cron.log